# target/linux/kirkwood/base-files/etc/hotplug.d/.svn/entries
10

dir
166
svn://184.82.232.92/myopenwrt/target/linux/kirkwood/base-files/etc/hotplug.d
svn://184.82.232.92/myopenwrt



2012-03-02T05:23:00.513110Z
5
zjhzzyf














e50a27de-cc27-4d45-b83a-890ea1c1ac35

usb
dir

block
dir

button
dir

